# cobbler/etc/supervisor/conf.d/supervisord.conf
[supervisord]
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
nodaemon=true
loglevel=info

[unix_http_server]
file=/var/run/supervisor.sock   ; (the path to the socket file)
chmod=0700                       ; sockef file mode (default 0700)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock ; use a unix:// URL  for a unix socket

[program:rsyslog]
command=/usr/sbin/rsyslogd -n
priority=1
#startsecs=5
redirect_stderr=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
autorestart=true
#autostart=false

[program:apache2]
command=apache2ctl -DFOREGROUND
priority=2
killasgroup=true
stopasgroup=true
autorestart=true
redirect_stderr=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
#autostart=false

[program:cobblerd]
command=/usr/bin/python /usr/local/bin/cobblerd -F
redirect_stderr=true
#autostart=false
priority=3
startsecs=5
redirect_stderr=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
autorestart=true
#autostart=false

[program:isc-dhcp-server]
command=/usr/sbin/dhcpd -4 -d -q -cf /etc/dhcp/dhcpd.conf eth0
priority=4
#istartsecs=10
redirect_stderr=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
autorestart=true
#autostart=false

#[program:in.tftpd]
#command=/usr/sbin/in.tftpd -4 -l -L -vvv -s /tftpboot
#priority=5
##startsecs=5
#redirect_stderr=true
#stdout_logfile=/dev/fd/1
#stdout_logfile_maxbytes=0
#autorestart=true
##autostart=false
[program:atftpd]
command=/usr/sbin/inetutils-inetd -d
priority=5
#startsecs=5
redirect_stderr=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
autorestart=true
#autostart=false
