# cobbler/tree/etc/supervisor/supervisord.conf
[supervisord]
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
nodaemon=true
loglevel=info

[unix_http_server]
file=/var/run/supervisor.sock   ; (the path to the socket file)
chmod=0700                       ; sockef file mode (default 0700)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock ; use a unix:// URL  for a unix socket

[program:update]
command=/usr/local/bin/update-config.sh
priority=1
redirect_stderr=true
autorestart=false


[program:rsyslog]
command=/usr/sbin/rsyslogd -n
priority=2
redirect_stderr=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
autorestart=true

[program:apache2]
command=apache2ctl -DFOREGROUND
priority=3
killasgroup=true
stopasgroup=true
autorestart=true
redirect_stderr=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0

[program:cobblerd]
command=/usr/bin/python3 /usr/local/bin/cobblerd -F
redirect_stderr=true
priority=4
startsecs=5
redirect_stderr=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
autorestart=true

[program:isc-dhcp-server]
command=/usr/sbin/dhcpd -4 -d -q -cf /etc/dhcp/dhcpd.conf eth0
priority=5
redirect_stderr=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
autorestart=true

[program:atftpd]
command=/usr/sbin/atftpd --daemon --tftpd-timeout 300 --retry-timeout 5 --no-multicast --no-fork --verbose=5 /srv/tftpboot/ --no-source-port-checking --logfile -
priority=6
redirect_stderr=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
autorestart=true
